# etc/cron.d/CRON_ddns_ufw_update
#
# DAT_CRON="/etc/cron.d/CRON_ddns_ufw_update" && echo "" > "${DAT_CRON}" && vi "${DAT_CRON}" && chown "root:root" "${DAT_CRON}" && chmod 0644 "${DAT_CRON}"; SERVICE_NAME="cron"; /bin/systemctl restart "${SERVICE_NAME}" 2>'/dev/null'; SERVICE_NAME="crond"; /bin/systemctl restart "${SERVICE_NAME}" 2>'/dev/null';
#

### [ DDNS ] Update DNS server(s) with current WAN IPv4
* * * * * root test -x "/usr/local/sbin/ddns_ufw_update" && "/usr/local/sbin/ddns_ufw_update";
