# etc/cron.d/CRON_jenkins_warfile_update
#
# DAT_CRON="/etc/cron.d/CRON_jenkins_warfile_update" && echo "" > "${DAT_CRON}" && vi "${DAT_CRON}" && chown "root:root" "${DAT_CRON}" && chmod 0644 "${DAT_CRON}" && service cron restart;
#

### [ JENKINS ] Update the main Jenkins executable (named a 'warfile', which has a '.war' file extension)
10 3 * * * root test -x "/usr/local/sbin/update_jenkins_warfile" && "/usr/local/sbin/update_jenkins_warfile";
