# etc/cron.d/CRON_jenkins_log_cleanup
#
# DAT_CRON="/etc/cron.d/CRON_jenkins_log_cleanup" && echo "" > "${DAT_CRON}" && vi "${DAT_CRON}" && chown "root:root" "${DAT_CRON}" && chmod 0644 "${DAT_CRON}" && service cron restart;
#

### [ JENKINS ] Cleanup excess logfiles created by Jenkins' runtime
25 01 * * * root test -x "/usr/local/sbin/jenkins_log_cleanup" && "/usr/local/sbin/jenkins_log_cleanup";
