# etc/cron.d/CRON_jenkins_war_update
#
# DAT_CRON="/etc/cron.d/CRON_update_jenkins_warfile" && echo "" > "${DAT_CRON}" && vi "${DAT_CRON}" && chown "root:root" "${DAT_CRON}" && chmod 0644 "${DAT_CRON}" && service cron restart;
#

10 3 * * * root /usr/local/sbin/update_jenkins_warfile;
