# etc/cron.d/jenkins_log_cleanup.cron
#
# vleanup jenkins' old build-logs
#

12 18 * * * sudo jenkins_log_cleanup;
