# etc/cron.d/jenkins_war_update.cron
#
# Update the Jenkins War-file (main executable)
#

10 3 * * * sudo jenkins_war_update;
