# etc/cron.d/jenkins_log_cleanup.cron
#
# vleanup jenkins' old build-logs
#

59 23 * * * root jenkins_log_cleanup;
