# etc/cron.d/update_jenkins_war
#
# Update the Jenkins War-file (main executable)
#

10 3 * * * update_jenkins_war;
